# roles/certbot/templates/certbot.timer.j2
# {{ ansible_managed }}

[Timer]
OnCalendar={{ certbot_system_timer_on_calender }} 
RandomizedDelaySec={{ certbot_system_timer_randomized_delay_sec }}
